# Sprint3/ServeurRasp/smartil.service
[Unit]
Description=Smartil service
Requires=network.target

[Service]
ExecStart=/home/Smartil/smartil.py